# install/certs/renew-cert.service
[Unit]
Description=Renew Certbot certificate
Wants=dhcpcd.service

[Service]
type=oneshot
ExecStart=/usr/bin/certbot renew
