# install/certs/renew-cert.timer
[Unit]
Description=Renew certificate

[Timer]
OnCalendar=daily
Persistant=true

[Install]
WantedBy=timers.target
